# root/etc/services.d/unbound/run
#!/usr/bin/with-contenv bash

chmod o+w /dev/stdout && \
exec \
	/usr/sbin/unbound -c /config/unbound.conf