# root/etc/services.d/unbound/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /usr/sbin/unbound -c /config/unbound.conf